# DOTRAIN 74/ping.conf
[program:ping]
process_name=%(program_name)s_%(process_num)02d
command= /usr/bin/python3 /home/viet/main.py
autostart=true
autorestart=true
user=viet
stdout_logfile=/home/viet/out1.log
